# template/supervisord.d/fluent.ini
[program:fluent]
user=fluent
command=/usr/local/bin/fluentd -vv -c /etc/fluent/fluent.conf
numprocs=1
autostart=true  ; supervisorが起動したら自動的に起動する
autorestart=true ; 落ちても自動的に再起動する
redirect_stderr=true ; エラー出力を標準出力にリダイレクトする
stdout_logfile=/var/log/supervisord/fluent.log
environment=LANG=ja_JP.UTF-8             ; 日本語環境に対応させる
