# template/supervisord.d/forecast.ini
[program:forecast]
user=forecast
command=/usr/local/bin/growthforecast.pl --data-dir /var/lib/forecast --front-proxy 127.0.0.1
numprocs=1
autostart=true  ; supervisorが起動したら自動的に起動する
autorestart=true ; 落ちても自動的に再起動する
redirect_stderr=true ; エラー出力を標準出力にリダイレクトする
stdout_logfile=/var/log/supervisord/forecast.log
