# template/supervisord.d/kahua.ini
[program:kahua]
user=kumano
command=/usr/local/bin/kahua-spvr -S /home/kumano/www/app
numprocs=1
autostart=true  ; supervisorが起動したら自動的に起動する
autorestart=true ; 落ちても自動的に再起動する
redirect_stderr=true ; エラー出力を標準出力にリダイレクトする
stdout_logfile=/var/log/supervisord/kahua.log
